# profile.d/keychain.sh
#! /bin/bash 
eval `keychain --quiet --eval`

